# xorg/.config/systemd/user/xresources.service
[Unit]
Description=Load Xresources

[Service]
Type=oneshot
ExecStart=/usr/sbin/xrdb %h/.config/xorg/Xresources
Restart=no
